# backend/fetch-emails.cron
# Run every 10 minutes (adjust as needed)
*/10 * * * * curl -s "http://localhost:7000/api/channels/fetch-emails?token=OrisRT2025" >> /var/log/cron.log 2>&1
